# berrycam_0.1-1_armhf/DEBIAN/prerm
#!/bin/bash
sudo systemctl stop berryCam
sudo systemctl disable berryCam
sudo systemctl daemon-reload
